# scripts/systemd/user/sort-home.path
[Unit]
Description=Watch home directory for changes

[Path]
PathModified=%h

[Install]
WantedBy=default.target
